# src/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
childlogdir=/var/log/
loglevel=info
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10


[unix_http_server]
file = /var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:pythonic]
command=/usr/local/bin/Pythonic
user=root
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stderr_logfile=/var/log/pythonic-stderr.log


[program:codeserver]
command=/usr/bin/code-server /root/Pythonic/executables
user=root
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stdout_logfile=/var/log/code-server-stdout.log
stderr_logfile=/var/log/code-server-stderr.log