# src/RPI/code-server.service
[Unit]
Description=Web based IDE

[Service]
Type=simple
User=pythonic
ExecStart=/usr/bin/code-server /home/pythonic/Pythonic/executables
